# systemd-services/web-server.service
[Unit]
Description=Development Web Server
Documentation=https://your-project-docs.com
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html
ExecStart=/usr/bin/python3 -m http.server 8080
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/html
PrivateTmp=true

# Environment
Environment=PYTHONPATH=/var/www/html
Environment=PORT=8080

[Install]
WantedBy=multi-user.target
